# packaging/DEBIAN/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create system user
        if ! id -u axon >/dev/null 2>&1; then
            useradd --system --no-create-home --shell /usr/sbin/nologin axon || true
        fi

        # Create required directories
        mkdir -p /var/lib/axon
        mkdir -p /run/axon
        mkdir -p /sys/fs/bpf/axon
        mkdir -p /etc/axon

        # Set permissions
        chown axon:axon /var/lib/axon 2>/dev/null || true
        chown axon:axon /run/axon 2>/dev/null || true
        chmod 750 /var/lib/axon

        # Mount bpffs if not mounted
        if ! mountpoint -q /sys/fs/bpf; then
            mount -t bpf bpffs /sys/fs/bpf 2>/dev/null || true
        fi

        # Install iptables rule for NFQUEUE steering (shared IP L7 traffic)
        iptables -t mangle -C PREROUTING -m mark --mark 0xBEEF -j NFQUEUE --queue-num 1 2>/dev/null || \
            iptables -t mangle -A PREROUTING -m mark --mark 0xBEEF -j NFQUEUE --queue-num 1 || true

        # Compile BPF object if clang is available
        if command -v clang >/dev/null 2>&1; then
            echo "Compiling eBPF XDP program..."
            clang -O2 -g -target bpf \
                -I/usr/include \
                -I/usr/include/$(uname -m)-linux-gnu \
                -D__KERNEL__ \
                -c /usr/lib/axon/xdp_firewall.c \
                -o /usr/lib/axon/xdp_firewall.o 2>/dev/null && \
                echo "eBPF compiled: /usr/lib/axon/xdp_firewall.o" || \
                echo "WARNING: eBPF compilation failed — daemon will run without XDP (L7 only)"
        else
            echo "WARNING: clang not found — eBPF/XDP unavailable"
            echo "         Install: apt install clang llvm linux-headers-\$(uname -r)"
        fi

        # Enable and start systemd service
        systemctl daemon-reload
        systemctl enable axon.service 2>/dev/null || true
        systemctl start axon.service 2>/dev/null || true

        echo ""
        echo "✅ axon installed successfully!"
        echo ""
        echo "   Service:     systemctl status axon"
        echo "   CLI help:    axon --help"
        echo "   Add rule:    axon add example.com --interface eth0"
        echo "   View logs:   axon logs -f"
        echo "   Log server:  localhost:5000 (JSON over TCP)"
        echo "   Change port: axon log-endpoint 0.0.0.0:9000"
        echo ""
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
